# etc/systemd/poppassd.socket
# /etc/systemd/system/poppassd.socket
[Unit]
Description=Poppassd Socket

[Socket]
ListenStream=106
Accept=yes

[Install]
WantedBy=sockets.target
